# harvest-actual/conf/etc/service/netapp-harvest/run
#!/bin/bash

exec /opt/netapp-harvest/netapp-worker -poller=harvest -conf=harvest-poller.conf
